# rootfs/etc/services.d/nut/run
#!/usr/bin/execlineb -P
s6-setuidgid nut
/usr/local/bin/python /opt/nut/nut.py -S
