# rootfs/etc/services.d/nut/run
#!/usr/bin/with-contenv bash

cd /nut

UMASK=${UMASK:-022}
umask "$UMASK"

s6-setuidgid nut python /nut/nut.py --scan --server